# app/supervisor/conf/001.conf
[program:mount]
command=/root/app/shell/mount.sh
autostart=false
startsecs=1
autorestart=true
startretries=3
user=root
priority=999
redirect_stderr=true
stdout_logfile=/root/app/shell/logs/mount.log
stopasgroup=false
killasgroup=false

[program:webdav]
command=/root/app/shell/webdav.sh
autostart=false
startsecs=1
autorestart=false
startretries=3
user=root
priority=999
redirect_stderr=true
stdout_logfile=/root/app/shell/logs/webdav.log
stopasgroup=true
killasgroup=true

[program:start]
command=/root/app/shell/start.sh
autostart=true
startsecs=1
autorestart=false
startretries=1
user=root
priority=999
redirect_stderr=true
stdout_logfile=/root/app/shell/logs/start.log
stopasgroup=true
killasgroup=true

[program:stop]
command=/root/app/shell/stop.sh
autostart=false
startsecs=1
autorestart=false
startretries=1
user=root
priority=999
redirect_stderr=true
stdout_logfile=/root/app/shell/logs/stop.log
stopasgroup=true
killasgroup=true